# skel/.config/autostart/osd-ip.desktop
[Desktop Entry]
Type=Application
Name=OSD-IP
Exec=/home/pi/ftcpi/bin/osd-ip
